# deploy/lumys.service
[Unit]
Description=Lumys OS Agent OS Daemon
Documentation=https://github.com/lumys-ai/lumys
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=lumys
Group=lumys
ExecStart=/usr/local/bin/lumys start
Restart=on-failure
RestartSec=5
TimeoutStopSec=30

# Environment
EnvironmentFile=-/etc/lumys/env
WorkingDirectory=/var/lib/lumys

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/lumys
PrivateTmp=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictSUIDSGID=true
MemoryDenyWriteExecute=false
RestrictRealtime=true

# Resource limits
LimitNOFILE=65536
LimitNPROC=4096

[Install]
WantedBy=multi-user.target
